# .config/systemd/user/keyboard-rate.service
[Unit]
Description=Set keyboard repeat rate for X11
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=oneshot
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
ExecStart=/usr/bin/xset r rate 250 40

[Install]
WantedBy=default.target
